# com.taronyang.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.taronyang.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zemyblue-mac-mini-m1/.nvm/versions/node/v22.23.0/bin/node</string>
        <string>dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/zemyblue-mac-mini-m1/.paperclip/instances/default/projects/d047328c-4364-4172-9730-3243342fcb49/c17cee8d-9cb0-4d6a-9838-d1d2150faac0/taronyang/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/taronyang-backend.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/taronyang-backend.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>FRONTEND_URL</key>
        <string>https://zemyalpha.github.io</string>
        <key>EXTRA_CORS_ORIGINS</key>
        <string>https://zemyalpha.github.io</string>
    </dict>
</dict>
</plist>
