# com.taronyang.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.taronyang.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zemyblue-mac-mini-m1/.paperclip/instances/default/projects/d047328c-4364-4172-9730-3243342fcb49/c17cee8d-9cb0-4d6a-9838-d1d2150faac0/taronyang/scripts/health-check.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/taronyang-monitor.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/taronyang-monitor.err</key>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/zemyblue-mac-mini-m1</string>
    </dict>
</dict>
</plist>
